# srcs/grafana/supervisor.conf
[supervisord]
nodaemon=true

[program:telegraf]
command = sh -c "telegraf && kill -s SIGTERM $(cat supervisord.pid)"

[program:grafana]
command = sh -c "grafana-server -homepath /usr/share/grafana -config /etc/default.ini && kill -s SIGTERM $(cat supervisord.pid)"
